# hp-vendor-opt-out.timer
[Unit]
Description=

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=4h

[Install]
WantedBy=timers.target
