# hp-vendor-daily.service
[Unit]
Description=
Documentation=

[Service]
Type=simple
ExecStart=/usr/libexec/hp-vendor daily
