# hp-vendor.service
[Unit]
Description=
Documentation=

[Service]
Type=simple
ExecStart=/usr/libexec/hp-vendor daemon

[Install]
WantedBy=graphical.target
